# debian/ntp-systemd-netif.service
[Service]
Environment=reason=BOUND
ExecStart=/bin/sh -c '. /etc/dhcp/dhclient-exit-hooks.d/ntp'

