# debian/ntp-systemd-netif.path
[Unit]
DefaultDependencies=no

[Path]
PathChanged=/run/systemd/netif/leases

[Install]
WantedBy=network-pre.target
