# debian/ntpsec.ntpsec-rotate-stats.service
[Unit]
Description=Rotate ntpd stats
Requisite=ntpsec.service

[Service]
Type=simple
ExecStart=/usr/libexec/ntpsec/rotate-stats
